# _unused/clamav/files/clamav-clamd.sh.in
#!/bin/sh
#
# PROVIDE: clamd
# REQUIRE: LOGIN
# BEFORE: mail 
# KEYWORD: shutdown

clamav_clamd_enable="${clamav_clamd_enable:=NO}"

clamav_clamd_socket="${clamav_clamd_socket:=%%CL_CLAMD_SOCKET%%}"

. /etc/rc.subr

name=clamav_clamd
rcvar=`set_rcvar`


command="%%PREFIX%%/sbin/clamd"
#pidfile="%%CL_RUN_DIR%%/clamd.pid"
required_dirs="%%CL_DB_DIR%%"
required_files="%%PREFIX%%/etc/clamav/clamd.conf"

load_rc_config "$name"

run_rc_command "$1"
#EOF
