# system/vm-bhyve/pkg-plist
etc/rc.d/vm
lib/vm-bhyve/vm-base
lib/vm-bhyve/vm-cmd
lib/vm-bhyve/vm-config
lib/vm-bhyve/vm-core
lib/vm-bhyve/vm-datastore
lib/vm-bhyve/vm-guest
lib/vm-bhyve/vm-info
lib/vm-bhyve/vm-migration
lib/vm-bhyve/vm-rctl
lib/vm-bhyve/vm-run
lib/vm-bhyve/vm-switch
lib/vm-bhyve/vm-switch-manual
lib/vm-bhyve/vm-switch-standard
lib/vm-bhyve/vm-switch-vale
lib/vm-bhyve/vm-switch-vxlan
lib/vm-bhyve/vm-util
lib/vm-bhyve/vm-zfs
man/man8/vm.8.gz
sbin/vm
share/examples/vm-bhyve/alpine.conf
share/examples/vm-bhyve/arch.conf
share/examples/vm-bhyve/centos6.conf
share/examples/vm-bhyve/centos7.conf
share/examples/vm-bhyve/config.sample
share/examples/vm-bhyve/coreos.conf
share/examples/vm-bhyve/debian.conf
share/examples/vm-bhyve/default.conf
share/examples/vm-bhyve/dragonfly.conf
share/examples/vm-bhyve/freebsd-zvol.conf
share/examples/vm-bhyve/freepbx.conf
share/examples/vm-bhyve/linux-zvol.conf
share/examples/vm-bhyve/netbsd.conf
share/examples/vm-bhyve/openbsd.conf
share/examples/vm-bhyve/resflash.conf
share/examples/vm-bhyve/ubuntu.conf
share/examples/vm-bhyve/windows.conf
@dir lib/vm-bhyve
@dir share/examples/vm-bhyve
@dir share/examples
