# my/vpnsw/pkg-plist
etc/rc.d/vpnsw
sbin/vpnsw
share/examples/vpnsw/vpnsw.conf.example
share/examples/vpnsw/vpnsw.crt.example
share/examples/vpnsw/vpnsw.key.example
share/examples/vpnsw/vpnsw.pw.example
share/vpnsw/public/css/app.css
share/vpnsw/public/css/datatables.css
share/vpnsw/public/css/datatables.min.css
share/vpnsw/public/css/foundation-float.css
share/vpnsw/public/css/foundation-float.min.css
share/vpnsw/public/css/foundation.css
share/vpnsw/public/css/foundation.min.css
share/vpnsw/public/favicon.ico
share/vpnsw/public/favicon.png
share/vpnsw/public/icons/foundation-icons.css
share/vpnsw/public/icons/foundation-icons.eot
share/vpnsw/public/icons/foundation-icons.svg
share/vpnsw/public/icons/foundation-icons.ttf
share/vpnsw/public/icons/foundation-icons.woff
share/vpnsw/public/icons/preview.html
share/vpnsw/public/images/sort_asc_disabled.png
share/vpnsw/public/images/sort_asc.png
share/vpnsw/public/images/sort_both.png
share/vpnsw/public/images/sort_desc_disabled.png
share/vpnsw/public/images/sort_desc.png
share/vpnsw/public/js/app.js
share/vpnsw/public/js/datatables.js
share/vpnsw/public/js/datatables.min.js
share/vpnsw/public/js/foundation.js
share/vpnsw/public/js/foundation.min.js
share/vpnsw/public/js/jquery.js
share/vpnsw/public/js/jquery.min.js
share/vpnsw/public/js/what-input.js
share/vpnsw/templs/exception.development.html.ep
share/vpnsw/templs/exception.production.html.ep
share/vpnsw/templs/hello.html.ep
share/vpnsw/templs/layouts/default.html.ep
share/vpnsw/templs/login.html.ep
share/vpnsw/templs/not_found.development.html.ep
share/vpnsw/templs/not_found.production.html.ep
share/vpnsw/templs/vpn-all.html.ep
share/vpnsw/templs/vpn-info.html.ep
share/vpnsw/templs/vpn-list.html.ep
@dir share/examples/vpnsw
@dir share/vpnsw/public/css
@dir share/vpnsw/public/icons
@dir share/vpnsw/public/images
@dir share/vpnsw/public/js
@dir share/vpnsw/public
@dir share/vpnsw/templs/layouts
@dir share/vpnsw/templs
@dir share/vpnsw
