# my/vam/pkg-plist
etc/rc.d/vam
sbin/vam
sbin/vam-helper
share/examples/vam/vam.conf.example
share/examples/vam/vam.crt.example
share/examples/vam/vam.key.example
share/examples/vam/vam.pw.example
share/vam/create-pg-db.sql
share/vam/public/css/app.css
share/vam/public/css/datatables.css
share/vam/public/css/datatables.min.css
share/vam/public/css/foundation-float.css
share/vam/public/css/foundation-float.min.css
share/vam/public/css/foundation.css
share/vam/public/css/foundation.min.css
share/vam/public/favicon.ico
share/vam/public/favicon.png
share/vam/public/icons/foundation-icons.css
share/vam/public/icons/foundation-icons.eot
share/vam/public/icons/foundation-icons.svg
share/vam/public/icons/foundation-icons.ttf
share/vam/public/icons/foundation-icons.woff
share/vam/public/icons/preview.html
share/vam/public/images/sort_asc_disabled.png
share/vam/public/images/sort_asc.png
share/vam/public/images/sort_both.png
share/vam/public/images/sort_desc_disabled.png
share/vam/public/images/sort_desc.png
share/vam/public/js/app.js
share/vam/public/js/datatables.js
share/vam/public/js/datatables.min.js
share/vam/public/js/foundation.js
share/vam/public/js/foundation.min.js
share/vam/public/js/jquery.js
share/vam/public/js/jquery.min.js
share/vam/public/js/what-input.js
share/vam/schema.sql
share/vam/templs/exception.development.html.ep
share/vam/templs/exception.production.html.ep
share/vam/templs/hello.html.ep
share/vam/templs/layouts/default.html.ep
share/vam/templs/login.html.ep
share/vam/templs/not_found.development.html.ep
share/vam/templs/not_found.production.html.ep
share/vam/templs/user-add-form.html.ep
share/vam/templs/user-add-handler.html.ep
share/vam/templs/user-delete-form.html.ep
share/vam/templs/user-delete-handler.html.ep
share/vam/templs/user-list.html.ep
share/vam/templs/user-rename-form.html.ep
share/vam/templs/user-rename-handler.html.ep
share/vam/templs/user-update-form.html.ep
share/vam/templs/user-update-handler.html.ep
@dir share/examples/vam
@dir share/vam/public/css
@dir share/vam/public/icons
@dir share/vam/public/images
@dir share/vam/public/js
@dir share/vam/public
@dir share/vam/templs/layouts
@dir share/vam/templs
@dir share/vam
