# net/frr8/files/isisd.in
#!/bin/sh
#
# PROVIDE: isisd
# REQUIRE: zebra
# KEYWORD: shutdown
#
. /etc/rc.subr

name="isisd"
rcvar="isisd_enable"
command="%%PREFIX%%/sbin/${name}"
pidfile="%%RUN_DIR%%/${name}.pid"
procname="${command}"
required_dirs="%%RUN_DIR%% %%LOG_DIR%%"
required_files="%%CONF_DIR%%/${name}.conf"
command_args=" -d"

load_rc_config $name

: ${isisd_enable="NO"}

run_rc_command "$1"
#EOF
