# net/frr8/files/ospfd.in
#!/bin/sh
#
# PROVIDE: ospfd
# REQUIRE: zebra
# KEYWORD: shutdown
#
. /etc/rc.subr

name="ospfd"
rcvar="ospfd_enable"
command="%%PREFIX%%/sbin/${name}"
pidfile="%%RUN_DIR%%/${name}.pid"
procname="${command}"
required_dirs="%%RUN_DIR%% %%LOG_DIR%%"
required_files="%%CONF_DIR%%/${name}.conf"
command_args=" -d"

load_rc_config $name

: ${ospfd_enable="NO"}

run_rc_command "$1"
#EOF
