# _beta/frr5/files/ripd.sh.in
#!/bin/sh
#
# PROVIDE: ripd
# REQUIRE: zebra
# KEYWORD: shutdown
#
. /etc/rc.subr

name="ripd"
rcvar="ripd_enable"
command="%%PREFIX%%/sbin/${name}"
#pidfile="%%RUN_DIR%%/${name}.pid"
procname="${command}"
required_dirs="%%RUN_DIR%% %%LOG_DIR%%"
required_files="%%CONF_DIR%%/${name}.conf"
command_args=" -d"

load_rc_config $name

: ${ripd_enable="NO"}

run_rc_command "$1"
#EOF
