# net/squid/files/_squid.sh.in
#!/bin/sh

# PROVIDE: squid
# REQUIRE: DAEMON

. /etc/rc.subr

name="squid"
rcvar="${name}_enable"
command="%%PREFIX%%/sbin/${name}"
pidfile="/var/run/${name}/${name}.pid"
procname="${name}"
required_files="%%PREFIX%%/etc/${name}/${name}.conf";
required_dirs="%%SQUID_SWAPDIR%% %%SQUID_LOGDIR%%";
sig_stop="KILL";
squid_user='%%SQUID_OWNER%%';


#stop_cmd=squid_stop_cmd
#
#squid_stop_cmd() {
#    echo "Stopping ${name}."
#    ${command} -k kill 2> /dev/null
#}

load_rc_config $name

: ${squid_enable="NO"}

run_rc_command "$1"
#EOF
