# data/mongodb60/pkg-plist
bin/mongod
bin/mongos
etc/rc.d/mongod
share/examples/mongodb/mongodb.conf.sample
@dir share/examples/mongodb
@dir share/examples
