# net/isc-dhcp/files/isc-dhcrelay6.in
#!/bin/sh
#
# PROVIDE: isc-dhcrelay6
# BEFORE: LOGIN
#

. /etc/rc.subr

name="isc_dhcrelay6"
rcvar="isc_dhcrelay6_enable"
command="%%PREFIX%%/sbin/isc-dhcrelay"
procname="${command}"
command_args="-6"

load_rc_config $name

: ${isc_dhcrelay6_enable="NO"}

run_rc_command "$1"
#EOF
