# net/isc-dhcp/files/isc-dhclient6.in
#!/bin/sh
#
# PROVIDE: isc_dhclient
# BEFORE: DAEMON
#
. /etc/rc.subr

name="isc_dhclient"
rcvar="isc_dhclient_enable"
command="%%PREFIX%%/sbin/isc-dhclient"
procname="${command}"
command_args=" -6"

load_rc_config $name

: ${isc_dhclient_enable="NO"}

run_rc_command "$1"
#EOF
