# net/openldap/files/slapd.sh.in
#!/bin/sh
#
# PROVIDE: slapd
# BEFORE: DAEMON
#
. /etc/rc.subr

name="slapd"
rcvar="${name}_enable"
command="%%PREFIX%%/sbin/${name}"
procname="${command}"
required_dirs="%%RUN_DIR%%"
required_files="%%CONF_DIR%%/${name}.conf"

load_rc_config $name

: ${slapd_flags="-h \"ldap:// ldaps:// ldapi://\""}
: ${slapd_enable="NO"}

run_rc_command "$1"
#EOF
