# net/dictd/files/dictd.sh.in
#!/bin/sh
#
# PROVIDE: dictd
# REQUIRE: DAEMON
#
# Add the following line to /etc/rc.conf[.local] to enable dictd:
#
# dictd_enable="YES"
#
. /etc/rc.subr

name=dictd
rcvar="dictd_enable"

command=%%PREFIX%%/sbin/dictd
required_files="%%PREFIX%%/etc/${name}.conf"

load_rc_config $name

: ${dictd_enable="NO"}

run_rc_command "$1"
#EOF
