# my/m2store/files/m2srv.in
#!/bin/sh

#
# PROVIDE: m2srv
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# m2srv_enable (bool):   Set to NO by default.
#               Set it to YES to enable m2srv.
#

. /etc/rc.subr

name="m2srv"
rcvar=m2srv_enable

command=%%PREFIX%%/sbin/${name}

load_rc_config $name

: ${m2srv_enable="NO"}

run_rc_command "$1"
