# my/m2store/files/s2srv.in
#!/bin/sh

#
# PROVIDE: s2srv
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# s2srv_enable (bool):   Set to NO by default.
#               Set it to YES to enable s2srv.
#

. /etc/rc.subr

name="s2srv"
rcvar=s2srv_enable

command=%%PREFIX%%/sbin/${name}

load_rc_config $name

: ${s2srv_enable="NO"}

run_rc_command "$1"
