# net/nfdump/files/nfcapd8.sh.in
#!/bin/sh
#
# PROVIDE: nfcapd8
# REQUIRE: DAEMON
#
. /etc/rc.subr

#instance=8

name="nfcapd8"
rcvar="nfcapd8_enable"
command="/usr/local/sbin/nfcapd"
pidfile="/var/run/nfcapd8.pid"
procname="${command}"
#command_args=" -D -P ${pidfile}"
required_dirs="/var/run"
#required_vars="${name}_interfaces ${name}_netflow_host ${name}_netflow_port"

load_rc_config $name

# /var/db/nfdump2 -S 8 -p 9998 -P /var/run/nfcapd2.pid

: ${nfcapd8_enable="NO"}
: ${nfcapd8_flags="-D -p 10003 -l /var/db/nfdump8 -P /var/run/nfcapd8.pid -S 8  "}


run_rc_command "$1"
#EOF
