# net/nfdump/files/nfcapd2.sh.in
#!/bin/sh
#
# PROVIDE: nfcapd2
# REQUIRE: DAEMON
#
. /etc/rc.subr

#instance=2

name="nfcapd2"
rcvar="nfcapd2_enable"
command="/usr/local/sbin/nfcapd"
pidfile="/var/run/nfcapd2.pid"
procname="${command}"
#command_args=" -D -P ${pidfile}"
required_dirs="/var/run"
#required_vars="${name}_interfaces ${name}_netflow_host ${name}_netflow_port"

load_rc_config $name

# /var/db/nfdump2 -S 8 -p 9998 -P /var/run/nfcapd2.pid

: ${nfcapd2_enable="NO"}
: ${nfcapd2_flags="-D -p 9997 -l /var/db/nfdump2 -P /var/run/nfcapd2.pid -S 8  "}


run_rc_command "$1"
#EOF
