# net/nfdump/files/nfcapd6.sh.in
#!/bin/sh
#
# PROVIDE: nfcapd6
# REQUIRE: DAEMON
#
. /etc/rc.subr

#instance=6

name="nfcapd6"
rcvar="nfcapd6_enable"
command="/usr/local/sbin/nfcapd"
pidfile="/var/run/nfcapd6.pid"
procname="${command}"
#command_args=" -D -P ${pidfile}"
required_dirs="/var/run"
#required_vars="${name}_interfaces ${name}_netflow_host ${name}_netflow_port"

load_rc_config $name

# /var/db/nfdump2 -S 8 -p 9998 -P /var/run/nfcapd2.pid

: ${nfcapd6_enable="NO"}
: ${nfcapd6_flags="-D -p 10001 -l /var/db/nfdump6 -P /var/run/nfcapd6.pid -S 8  "}


run_rc_command "$1"
#EOF
