# net/nfdump/files/nfcapd5.sh.in
#!/bin/sh
#
# PROVIDE: nfcapd5
# REQUIRE: DAEMON
#
. /etc/rc.subr

#instance=5

name="nfcapd5"
rcvar="nfcapd5_enable"
command="/usr/local/sbin/nfcapd"
pidfile="/var/run/nfcapd5.pid"
procname="${command}"
#command_args=" -D -P ${pidfile}"
required_dirs="/var/run"
#required_vars="${name}_interfaces ${name}_netflow_host ${name}_netflow_port"

load_rc_config $name

# /var/db/nfdump2 -S 8 -p 9998 -P /var/run/nfcapd2.pid

: ${nfcapd5_enable="NO"}
: ${nfcapd5_flags="-D -p 10000 -l /var/db/nfdump5 -P /var/run/nfcapd5.pid -S 8  "}


run_rc_command "$1"
#EOF
