# net/nfdump/files/nfcapd7.sh.in
#!/bin/sh
#
# PROVIDE: nfcapd7
# REQUIRE: DAEMON
#
. /etc/rc.subr

#instance=7

name="nfcapd7"
rcvar="nfcapd7_enable"
command="/usr/local/sbin/nfcapd"
pidfile="/var/run/nfcapd7.pid"
procname="${command}"
#command_args=" -D -P ${pidfile}"
required_dirs="/var/run"
#required_vars="${name}_interfaces ${name}_netflow_host ${name}_netflow_port"

load_rc_config $name

# /var/db/nfdump2 -S 8 -p 9998 -P /var/run/nfcapd2.pid

: ${nfcapd7_enable="NO"}
: ${nfcapd7_flags="-D -p 10002 -l /var/db/nfdump7 -P /var/run/nfcapd7.pid -S 8  "}


run_rc_command "$1"
#EOF
