# net/nfdump/files/nfcapd4.sh.in
#!/bin/sh
#
# PROVIDE: nfcapd4
# REQUIRE: DAEMON
#
. /etc/rc.subr

#instance=4

name="nfcapd4"
rcvar="nfcapd4_enable"
command="/usr/local/sbin/nfcapd"
pidfile="/var/run/nfcapd4.pid"
procname="${command}"
#command_args=" -D -P ${pidfile}"
required_dirs="/var/run"
#required_vars="${name}_interfaces ${name}_netflow_host ${name}_netflow_port"

load_rc_config $name

# /var/db/nfdump2 -S 8 -p 9998 -P /var/run/nfcapd2.pid

: ${nfcapd4_enable="NO"}
: ${nfcapd4_flags="-D -p 9999 -l /var/db/nfdump4 -P /var/run/nfcapd4.pid -S 8  "}


run_rc_command "$1"
#EOF
