# _unused/nfdump/files/sfcapd.sh.in
#!/bin/sh
#
# PROVIDE: sfcapd
# REQUIRE: DAEMON
#
. /etc/rc.subr

name="sfcapd"
rcvar="sfcapd_enable"
command="%%PREFIX%%/sbin/${name}"
pidfile="%%RUN_DIR%%/${name}.pid"
procname="${command}"
command_args=" -D -P ${pidfile}"

load_rc_config $name

: ${sfcapd_enable="NO"}
: ${sfcapd_datadir="/var/${name}"}
: ${sfcapd_port=6343}
: ${sfcapd_schema=8}
: ${sfcapd_time=120}
: ${sfcapd_args=""}

required_dirs="%%RUN_DIR%% ${sfcapd_datadir}"

: ${sfcapd_flags="-l ${sfcapd_datadir} -t ${sfcapd_time} -p ${sfcapd_port} -S ${sfcapd_schema} ${sfcapd_args}"}

run_rc_command "$1"
#EOF
